# api/runit/wemo3000/run
#!/bin/sh
export WEMO3000_RELEASE=1
export GIN_MODE=release
exec 2>&1
exec wemo3000-server
